# .config/systemd/user/pgbackup.service
[Unit]
Description=PostgreSQL Backup for Syncthing

[Service]
Type=oneshot
ExecStart=/home/saltchicken/.config/systemd/user/scripts/pgbackup.sh
